# config/systemd/user/swayosd-keepalive.service
[Unit]
Description=Keep swayosd warm

[Service]
Type=oneshot
ExecStart=hyprkarl-swayosd-poke
